# etc/sv/mclient/finish
#!/bin/bash
pkill omxplayer
fbset -depth 8
fbset -depth 16
umount /dev/sda

